# smf/conch.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
        Created by Manifold
--><service_bundle type="manifest" name="conch">

    <service name="site/conch" type="service" version="1">

        <create_default_instance enabled="false"/>

        <single_instance/>

        <dependency name="network" grouping="require_all" restart_on="error" type="service">
            <service_fmri value="svc:/milestone/network:default"/>
        </dependency>

        <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
            <service_fmri value="svc:/system/filesystem/local"/>
        </dependency>


        <method_context>

        </method_context>
<exec_method type="method" name="start" exec="/opt/local/lib/perl5/site_perl/bin/carton exec hypnotoad bin/conch" timeout_seconds="60">
            <method_context working_directory='/opt/conch/Conch'></method_context>
        </exec_method>

        <exec_method type="method" name="refresh" exec="/opt/local/lib/perl5/site_perl/bin/carton exec hypnotoad bin/conch" timeout_seconds="60">
            <method_context working_directory='/opt/conch/Conch'></method_context>
        </exec_method>

        <exec_method type="method" name="stop" exec="/opt/local/lib/perl5/site_perl/bin/carton exec hypnotoad -s bin/conch" timeout_seconds="60">
            <method_context working_directory='/opt/conch/Conch'></method_context>
        </exec_method>

        <property_group name="startd" type="framework">


            <propval name="duration" type="astring" value="contract"/>
            <propval name="ignore_error" type="astring" value="core,signal"/>
        </property_group>

        <property_group name="application" type="application">

        </property_group>


        <stability value="Evolving"/>

        <template>
            <common_name>
                <loctext xml:lang="C">
                    Conch API Service
                </loctext>
            </common_name>
        </template>

    </service>

</service_bundle>

